# templates/act_runner.service.j2
[Unit]
Description=Gitea Actions Runner
After=network.target gitea.service
Wants=gitea.service

[Service]
Type=simple
User=act_runner
Group=act_runner
WorkingDirectory=/var/lib/act_runner
ExecStart=/usr/local/bin/act_runner daemon --config /etc/act_runner/config.yaml
Restart=always
RestartSec=2
Environment=HOME=/var/lib/act_runner

[Install]
WantedBy=multi-user.target
